# FV3GFSwfm/rt_ufscam_l64/rt_ufscam_l64.crontab

#################### rt_ufscam_l64 ####################
MAILTO=""
*/5 * * * * /apps/rocoto/1.3.1/bin/rocotorun -d /lfs4/HFIP/gsd-fv3-hfip/rtruns/UFS-CAMsuite//FV3GFSwfm/rt_ufscam_l64/rt_ufscam_l64.db -w /lfs4/HFIP/gsd-fv3-hfip/rtruns/UFS-CAMsuite//FV3GFSwfm/rt_ufscam_l64/rt_ufscam_l64.xml
#################################################################

